# services.d/test/run
#!/usr/bin/execlineb -P

exec php /signal/signal.php