# services.d/test/run
#!/bin/bash

function handler {
	echo "$(date) Caught signal" | tee -a /tmp/signal.txt
	exit
}

trap handler SIGINT SIGTERM

exec sleep 999999
