# services.d/test/run
#!/bin/bash

function handler {
	echo "$(date) Caught signal" | tee -a /tmp/signal.txt
}

trap handler SIGINT SIGTERM

sleep 999999
